# dist/systemd/djinn-curator.service


[Unit]
Description=Djinn CI Curator
After=network-online.target

[Service]
Type=simple
User=djinn-curator
Group=djinn
PIDFile=/var/run/djinn/curator.pid
ExecStart=/usr/bin/djinn-curator -config /etc/djinn/curator.conf
KillMode=mixed
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
