# dist/systemd/djinn-server.service


[Unit]
Description=Djinn CI Server
After=network-online.target

[Service]
Type=simple
User=djinn-server
Group=djinn
PIDFile=/var/run/djinn/server.pid
ExecStart=/usr/bin/djinn-server -config /etc/djinn/server.conf
KillMode=mixed
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
