# dist/systemd/djinn-worker.service


[Unit]
Description=Djinn CI Worker
After=network-online.target

[Service]
Type=simple
User=djinn-worker
Group=djinn
PIDFile=/var/run/djinn/worker.pid
ExecStart=/usr/bin/djinn-worker -config /etc/djinn/worker.conf -driver /etc/djinn/driver.conf
KillMode=mixed
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
